# supervisor.ini
[program:ecommerce]
command=/data/venvs/ecommerce/bin/gunicorn --chdir /data/projects/ecommerce-backend ecommerce.wsgi:application --reload --workers=2 -b 0.0.0.0:7079 --access-logfile gunicorn-access.log --error-logfile gunicorn-error.log
directory= /data/projects/ecommerce-backend
stdout_logfile=/data/logs/ecommerce-backend/gunicorn-stdout.log
stderr_logfile=/data/code/ecommerce-backend/gunicorn-stderr.log
autostart=true
autorestart=true

[program:ecommerce_celery_worker]
command=/data/venvs/ecommerce/bin/celery ecommerce.celery worker  -A -l info --concurrency=2 --logfile=/data/logs/ecommerce-backend/celery_worker.log
user=celery
directory= /data/projects/ecommerce-backend
stdout_logfile=/data/logs/ecommerce-backend/celery_worker-stdout.log
stderr_logfile=/data/code/ecommerce-backend/celery_worker-stderr.log
autostart=true
autorestart=true

[program:ecommerce_celery_beat]
command=/data/venv/ecommerce/bin/celery ecommerce.celery beat -A -l info --logfile=/data/logs/ecommerce-backend/celery_beat.log
user=celery
directory=/data/code/ecommerce-backend
stdout_logfile=/data/logs/ecommerce-backend/celery_beat-stdout.log
stderr_logfile=/data/logs/ecommerce-backend/celery_beat-stderr.log
